# deploy/llmexpensetracker.service
[Unit]
Description=LLM Expense Tracker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=llmexpense
WorkingDirectory=/opt/LlmExpenseTracker
EnvironmentFile=/opt/LlmExpenseTracker/.env
ExecStart=/opt/LlmExpenseTracker/bin/server
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
